# noblinks-agent/templates/noblinks-agent.service
[Unit]
Description=Noblinks Agent
Documentation=https://noblinks.io/docs/agent
After=network-online.target prometheus.service
Wants=network-online.target
Requires=prometheus.service

[Service]
Type=simple
User=root
EnvironmentFile=/etc/noblinks/agent.env
ExecStart=/usr/local/bin/noblinks-agent
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noblinks-agent

[Install]
WantedBy=multi-user.target
